# package/etc/init/canary-agent.conf
start on startup 
stop on shutdown
 
respawn
respawn limit 10 5
console log

script
  exec canary-agent
end script
