# package/ubuntu/trusty/files/etc/init/appcanary.conf
start on filesystem
stop on runlevel [06]
 
respawn
respawn limit 10 5
console log

script
  exec appcanary
end script
