# scripts/debian/jessie/etc/systemd/system/appcanary.service
[Unit]
Description=appcanary agent
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/appcanary

[Install]
WantedBy=multi-user.target
